# scripts/deployment/data/supervisord.conf
[supervisord]
nodaemon=true
logfile=/lerna/build/logs/supervisor.log
pidfile=/lerna/build/supervisor.pid

[program:lerna_gunicorn]
command=gunicorn lerna.wsgi:application --workers 3 --bind unix:/lerna/build/gunicorn.sock --log-file /lerna/build/logs/gunicorn.log
directory=/lerna
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10

